# prebuildfs/etc/profile.d/common-bin.sh
# Add ~/common/bin to the PATH
pathmunge ${HOME}/common/bin after
